# com.birdblues.investmentlog.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.birdblues.investmentlog.daemon</string>
    <key>UserName</key>
    <string>birdblues</string>
    <key>GroupName</key>
    <string>staff</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>/Users/birdblues/workspace/InvestmentLog/main_local.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/birdblues/workspace/InvestmentLog</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/birdblues</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/birdblues/workspace/InvestmentLog/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/birdblues/workspace/InvestmentLog/logs/stderr.log</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
